# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/etc/init.d/nginx restart

[program:mysqld]
command=/usr/local/bin/run
